# deploy/vca.service
# systemd 服务单元：让 VCA 在服务器上常驻、开机自启、崩溃自动重启
# 安装：
#   1) 把 jar 放到 /opt/vca/app.jar
#   2) 把真实 key 写进 /etc/vca.env（见下方 EnvironmentFile 示例）
#   3) sudo cp deploy/vca.service /etc/systemd/system/vca.service
#   4) sudo systemctl daemon-reload && sudo systemctl enable --now vca
#   5) 查看日志：journalctl -u vca -f

[Unit]
Description=VCA Voice Chat Assistant
After=network.target

[Service]
Type=simple
# 建议用专用非 root 用户运行
User=vca
WorkingDirectory=/opt/vca
EnvironmentFile=/etc/vca.env
# 2核2G 服务器：堆控制在 512m，给 OS/Caddy/Netty直接内存 留余量
ExecStart=/usr/bin/java -Xms256m -Xmx512m -XX:MaxDirectMemorySize=128m -Duser.timezone=Asia/Shanghai -jar /opt/vca/app.jar
SuccessExitStatus=143
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target

# ---- /etc/vca.env 示例（自行创建，权限 600）----
# SPRING_PROFILES_ACTIVE=real
# SERVER_PORT=8080
# DASHSCOPE_API_KEY=sk-你的阿里云key
# DEEPSEEK_ENABLED=true
# DEEPSEEK_API_KEY=sk-你的deepseek-key
# DEEPSEEK_PROXY=
